# branches/LOFAR-Task8898/SAS/OTDB_Services/TreeService.ini
[program:TreeService]
command=/bin/bash -c 'source %(ENV_LOFARROOT)s/lofarinit.sh;TreeService.py'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stdout_logfile=%{ENV_LOFARROOT}s/var/log/%{program_name)s.stdout
stderr_logfile=%{ENV_LOFARROOT}s/var/log/%{program_name)s.stderr
