# branches/LOFAR-Task8898/SAS/OTDB_Services/TreeStatusEvents.ini
[program:TreeStatusEvents]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;TreeStatusEvents.py --hostname=sasdb --database=LOFAR_4 --busname lofar.otdb.status'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stdout_logfile=%{program_name)s.stdout
stderr_logfile=%{program_name)s.stderr
