# branches/RA-Task8887/SAS/ResourceAssignment/SystemStatusService/ssdbservice.ini
[program:ssdbservice]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;ssdbservice --busname=lofar.system --servicename=SSDBService'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stopasgroup=true ; bash does not propagate signals
stdout_logfile=%(program_name)s.log
redirect_stderr=true
stderr_logfile=NONE
