# MAC/Services/src/pipelinecontrol.ini
[program:PipelineStarter]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;pipelinecontrol'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stopasgroup=true
stdout_logfile=%(program_name)s.log
stderr_logfile=%(program_name)s.stderr
