# SAS/ResourceAssignment/Services/src/jobstoscheduleservice.ini
[program:JobsToSchedule]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;jobstoscheduleservice --otdb_busname=lofar.otdb --my_busname=lofar.ra'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stopasgroup=true
stdout_logfile=%(program_name)s.log
stderr_logfile=%(program_name)s.stderr
