# branches/LOFAR-Task8888/SAS/ResourceAssignment/ResourceAssignmentEditor/etc/resourceassignmenteditorwebservice.ini
[program:ResourceAssignmentEditor]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;resourceassignmenteditorwebservice.py'
user=nginx
stopsignal=INT ; KeyboardInterrupt
stopasgroup=true
stdout_logfile=%(program_name)s.log
stderr_logfile=%(program_name)s.stderr
