# branches/LOFAR-SSDB-Task9026/SAS/MoM/MoMQueryService/momqueryservice.ini
[program:momqueryservice]
command=/bin/bash -c 'source $LOFARROOT/lofarinit.sh;momqueryservice --busname=lofar.ra.command --servicename=momqueryservice'
user=lofarsys
stopsignal=INT ; KeyboardInterrupt
stopasgroup=true ; bash does not propagate signals
stdout_logfile=%(program_name)s.log
redirect_stderr=true
stderr_logfile=NONE
